# roles/managed/templates/ansible-pull.timer
[Unit]
Description=Run Ansible pull periodically

[Timer]
OnBootSec=2min
OnUnitActiveSec=1h

[Install]
WantedBy=timers.target
